# examples/systemd/ruyi-backend-process-telemetry.timer
[Unit]
Description=RuyiSDK Backend: Process telemetry data (hourly)

[Timer]
OnCalendar=*:0

[Install]
WantedBy=timers.target
